# vastai-images/desktop-editor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:vnc]
command=/usr/local/bin/vnc-startup.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/vnc.log
stderr_logfile=/var/log/vnc.err

[program:novnc]
command=/usr/share/novnc/utils/launch.sh --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
